# deploy/systemd/private-example/mtproto/tracegate-mtproto@.service.example
[Unit]
Description=Tracegate private MTProto gateway (%i)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/tracegate/private/mtproto/run-mtproto.sh

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/tracegate/private/mtproto/mtproto.env
Environment=CONFIG_DIR=/etc/tracegate
Environment=TRACEGATE_RUNTIME_ROLE=%i
ExecStart=/usr/bin/env bash /etc/tracegate/private/mtproto/run-mtproto.sh start %i
ExecReload=/usr/bin/env bash /etc/tracegate/private/mtproto/run-mtproto.sh reload %i
ExecStop=/usr/bin/env bash /etc/tracegate/private/mtproto/run-mtproto.sh stop %i

[Install]
WantedBy=multi-user.target
