# deploy/systemd/private-example/systemd/tracegate-obfuscation@.service.example
[Unit]
Description=Tracegate private obfuscation helper (%i)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/tracegate/private/systemd/run-obfuscation.sh

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/tracegate/private/systemd/obfuscation.env
Environment=CONFIG_DIR=/etc/tracegate
Environment=TRACEGATE_RUNTIME_ROLE=%i
ExecStart=/usr/bin/env bash /etc/tracegate/private/systemd/run-obfuscation.sh start %i
ExecReload=/usr/bin/env bash /etc/tracegate/private/systemd/run-obfuscation.sh reload %i
ExecStop=/usr/bin/env bash /etc/tracegate/private/systemd/run-obfuscation.sh stop %i
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
