# deploy/systemd/tracegate-haproxy@.service
[Unit]
Description=Tracegate HAProxy fronting runtime (%i)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/sbin/haproxy
ConditionPathExists=/var/lib/tracegate/agent-%i/runtime/haproxy/haproxy.cfg

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/lib/tracegate/agent-%i/runtime/haproxy
ExecStart=/usr/sbin/haproxy -W -db -f /var/lib/tracegate/agent-%i/runtime/haproxy/haproxy.cfg
Restart=on-failure
RestartSec=2
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/tracegate/agent-%i/runtime /var/log/tracegate /run
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
