# deploy/systemd/tracegate-agent-transit.service
[Unit]
Description=Tracegate agent (Transit)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/tracegate
EnvironmentFile=/etc/tracegate/tracegate.env
EnvironmentFile=-/etc/tracegate/transit.env
Environment=AGENT_RUNTIME_MODE=systemd
Environment=AGENT_ROLE=TRANSIT
ExecStart=/opt/tracegate/.venv/bin/tracegate-agent
Restart=on-failure
RestartSec=3
UMask=0077
StateDirectory=tracegate
ConfigurationDirectory=tracegate
RuntimeDirectory=tracegate
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
