# deploy/systemd/tracegate-hysteria@.service
[Unit]
Description=Tracegate Hysteria2 UDP runtime (%i)
After=network-online.target tracegate-agent-%i.service
Wants=network-online.target
ConditionPathExists=/usr/local/bin/hysteria
ConditionPathExists=/var/lib/tracegate/agent-%i/runtime/hysteria/server.yaml

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/lib/tracegate/agent-%i/runtime/hysteria
ExecStart=/usr/local/bin/hysteria server -c /var/lib/tracegate/agent-%i/runtime/hysteria/server.yaml
Restart=on-failure
RestartSec=2
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/tracegate/agent-%i/runtime /var/log/tracegate
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
