# deploy/systemd/tracegate-api.service
[Unit]
Description=Tracegate API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tracegate
Group=tracegate
WorkingDirectory=/opt/tracegate
EnvironmentFile=/etc/tracegate/control-plane.env
ExecStart=/opt/tracegate/.venv/bin/tracegate-api
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
