# deploy/systemd/tracegate-nginx@.service
[Unit]
Description=Tracegate Nginx TLS / decoy runtime (%i)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/sbin/nginx
ConditionPathExists=/var/lib/tracegate/agent-%i/runtime/nginx/nginx.conf

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/lib/tracegate/agent-%i/runtime/nginx
ExecStart=/usr/sbin/nginx -g 'daemon off;' -c /var/lib/tracegate/agent-%i/runtime/nginx/nginx.conf
ExecReload=/usr/sbin/nginx -s reload -c /var/lib/tracegate/agent-%i/runtime/nginx/nginx.conf
Restart=on-failure
RestartSec=2
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/tracegate/agent-%i/runtime /var/log/tracegate /run /etc/tracegate
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
