# deploy/systemd/tracegate-agent.service
[Unit]
Description=Tracegate Node Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tracegate
Group=tracegate
WorkingDirectory=/opt/tracegate
EnvironmentFile=/etc/tracegate/agent.env
ExecStart=/opt/tracegate/.venv/bin/tracegate-agent
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
