# deploy/systemd/tracegate-wireguard.service
[Unit]
Description=Tracegate WireGuard server
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/wg-quick up /etc/tracegate/private/wireguard/wg.conf
ExecStop=/usr/bin/wg-quick down /etc/tracegate/private/wireguard/wg.conf

[Install]
WantedBy=multi-user.target
