# deploy/systemd/tracegate-dispatcher.service
[Unit]
Description=Tracegate dispatcher
After=network-online.target tracegate-api.service
Wants=network-online.target

[Service]
Type=simple
User=tracegate
Group=tracegate
WorkingDirectory=/opt/tracegate
EnvironmentFile=/etc/tracegate/tracegate.env
ExecStart=/opt/tracegate/.venv/bin/tracegate-dispatcher
Restart=on-failure
RestartSec=3
UMask=0077
StateDirectory=tracegate
ConfigurationDirectory=tracegate
RuntimeDirectory=tracegate
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
