# deploy/systemd/private-example/fronting/tracegate-fronting@.service.example
[Unit]
Description=Tracegate private TCP/443 fronting helper (%i)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/tracegate/private/fronting/run-fronting.sh

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/tracegate/private/fronting/fronting.env
Environment=CONFIG_DIR=/etc/tracegate
Environment=TRACEGATE_RUNTIME_ROLE=%i
ExecStart=/usr/bin/env bash /etc/tracegate/private/fronting/run-fronting.sh start %i
ExecReload=/usr/bin/env bash /etc/tracegate/private/fronting/run-fronting.sh reload %i
ExecStop=/usr/bin/env bash /etc/tracegate/private/fronting/run-fronting.sh stop %i

[Install]
WantedBy=multi-user.target
