# deploy/systemd/private-example/profiles/tracegate-profiles@.service.example
[Unit]
Description=Tracegate private profile adapter helper (%i)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/tracegate/private/profiles/run-profiles.sh

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/tracegate/private/profiles/profiles.env
Environment=CONFIG_DIR=/etc/tracegate
Environment=TRACEGATE_RUNTIME_ROLE=%i
ExecStart=/usr/bin/env bash /etc/tracegate/private/profiles/run-profiles.sh start %i
ExecReload=/usr/bin/env bash /etc/tracegate/private/profiles/run-profiles.sh reload %i
ExecStop=/usr/bin/env bash /etc/tracegate/private/profiles/run-profiles.sh stop %i

[Install]
WantedBy=multi-user.target
