# deploy/systemd/private-example/link-crypto/tracegate-link-crypto@.service.example
[Unit]
Description=Tracegate private link-crypto helper (%i)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/tracegate/private/link-crypto/run-link-crypto.sh

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/tracegate/private/link-crypto/link-crypto.env
Environment=CONFIG_DIR=/etc/tracegate
Environment=TRACEGATE_RUNTIME_ROLE=%i
ExecStart=/usr/bin/env bash /etc/tracegate/private/link-crypto/run-link-crypto.sh start %i
ExecReload=/usr/bin/env bash /etc/tracegate/private/link-crypto/run-link-crypto.sh reload %i
ExecStop=/usr/bin/env bash /etc/tracegate/private/link-crypto/run-link-crypto.sh stop %i

[Install]
WantedBy=multi-user.target
